# roles/mail/files/mbsync.timer
[Unit]
Description=Mailbox sync timer

[Timer]
OnBootSec=10
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
