# roles/mail/files/notmuch_new.service
[Unit]
Description=Index new mail

[Service]
Type=oneshot
ExecStart=/usr/bin/notmuch new

[Install]
WantedBy=default.target
