# roles/mail/files/notmuch_new.timer
[Unit]
Description=Mail indexer timer

[Timer]
OnBootSec=10
OnUnitActiveSec=7min
Persistent=true

[Install]
WantedBy=timers.target
